# supabase/functions/update-expired-dealerships/crontab
0 0 * * * /bin/bash -c 'curl -i -X POST https://momhdxlmwzdikxmwittx.functions.supabase.co/update-expired-dealerships'